# test/e2e/quadlet/image-volume.container
## assert-podman-final-args localhost/imagename
## assert-podman-args --image-volume tmpfs

[Container]
Image=localhost/imagename
ImageVolume=tmpfs
